# install/linux/dmx-server-manager.service
# Managed by the DmxServerManager installer.
[Unit]
Description=DmxServerManager game server manager
Documentation=https://github.com/thefrcrazy/DmxServerManager
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=dmx-server-manager
Group=dmx-server-manager
WorkingDirectory=/usr/lib/dmx-server-manager/current
Environment=DMX_CONFIG_FILE=/etc/dmx-server-manager/config.toml
Environment=DMX_DATA_DIR=/var/lib/dmx-server-manager
Environment=DMX_MASTER_KEY_FILE=/etc/dmx-server-manager/master.key
Environment=DMX_STEAMCMD_PATH=/usr/games/steamcmd
Environment=DMX_STATIC_DIR=/usr/lib/dmx-server-manager/current/static
Environment=DMX_DEPLOYMENT_MODE=native
Environment=HOME=/var/lib/dmx-server-manager
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/usr/lib/dmx-server-manager/current/dmx-server-manager
Restart=on-failure
RestartSec=5s
TimeoutStartSec=90s
TimeoutStopSec=120s
KillMode=control-group
UMask=0077
AmbientCapabilities=
CapabilityBoundingSet=
NoNewPrivileges=true
PrivateTmp=true
ProtectClock=true
ProtectControlGroups=true
ProtectHome=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectProc=invisible
ProtectSystem=strict
ReadWritePaths=/var/lib/dmx-server-manager
RemoveIPC=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictSUIDSGID=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
